# runeio/generic/Packages-Dir/package-rune-etc/hotplug.d/ntp/ntp_time_adjusted.sh
#!/bin/sh

/bin/echo -n "NTP Hotplug Time adjusted action ${1} - ${ACTION} at " > /tmp/ntp_time_adjusted.log
/bin/date >> /tmp/ntp_time_adjusted.log
exit 0
